# systemd/unitc-external-model-highs.service
[Unit]
Description=UnitC external model executor (HiGHS piecewise MILP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Coding/unitc/backend
EnvironmentFile=/srv/unitc-external-model/model.env
ExecStart=%h/.venvs/unitc-external-model/bin/python -m app.serve --host 127.0.0.1 --port 18003 --model-id highs_piecewise_milp --model-label "Reference Internal (HiGHS)"
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
